# backend/supervisor/celery.conf
; Supervisor configuration for Celery
; Copy to /etc/supervisor/conf.d/celery.conf
; Then run: supervisorctl reread && supervisorctl update

[program:celery-worker]
command=/app/venv/bin/celery -A app.core.celery_app worker --loglevel=info --pool=prefork --concurrency=4 --hostname=worker-%(process_num)02d@%%h
directory=/app
user=www-data
numprocs=2
process_name=%(program_name)s-%(process_num)02d
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stopasgroup=true
killasgroup=true
priority=998
stdout_logfile=/var/log/celery/worker-%(process_num)02d.log
stderr_logfile=/var/log/celery/worker-%(process_num)02d-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

; Environment variables
environment=
    PYTHONPATH="/app",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s"

[program:celery-beat]
command=/app/venv/bin/celery -A app.core.celery_app beat --loglevel=info --schedule=/tmp/celerybeat-schedule
directory=/app
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=999
stdout_logfile=/var/log/celery/beat.log
stderr_logfile=/var/log/celery/beat-error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

; Environment variables
environment=
    PYTHONPATH="/app",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s"

[group:celery]
programs=celery-worker,celery-beat
priority=999
